# rootfs/pkgs/pppox/filelist
etc/apm/event.d/ppp
etc/bash_completion.d/pon
etc/chatscripts/pap
etc/init.d/pppd-dns
etc/ppp
etc/ppp/pppoe-server-options
etc/ppp/ppp_on_boot.dsl
etc/ppp/ip-down
etc/ppp/ipv6-up
etc/ppp/pap-secrets
etc/ppp/options
etc/ppp/chap-secrets
etc/ppp/ip-down.d
etc/ppp/ip-down.d/0000usepeerdns
etc/ppp/peers
etc/ppp/peers/dsl-provider
etc/ppp/radius
etc/ppp/radius/servers
etc/ppp/radius/dictionary.asnet
etc/ppp/radius/dictionary
etc/ppp/radius/dictionary.microsoft
etc/ppp/ip-up
etc/ppp/ip-up.d
etc/ppp/ip-up.d/0000usepeerdns
etc/ppp/ipv6-down
usr/lib/pppd
usr/lib/pppd/2.4.4
usr/lib/pppd/2.4.4/rp-pppoe.so
usr/lib/pppd/2.4.4/radattr.so
usr/lib/pppd/2.4.4/radrealms.so
usr/lib/pppd/2.4.4/passwordfd.so
usr/lib/pppd/2.4.4/radius.so
usr/lib/pppd/2.4.4/minconn.so
usr/lib/pppd/2.4.4/pppoatm.so
usr/lib/pppd/2.4.4/winbind.so
usr/lib/pppd/2.4.4/passprompt.so
usr/sbin/chat
usr/sbin/pppd
usr/sbin/pppdump
usr/sbin/pppoe
usr/sbin/pppoe-connect
usr/sbin/pppoe-discovery
usr/sbin/pppoe-relay
usr/sbin/pppoe-server
usr/sbin/pppoe-setup
usr/sbin/pppoe-sniff
usr/sbin/pppoe-start
usr/sbin/pppoe-status
usr/sbin/pppoe-stop
usr/sbin/pppstats
usr/share/ppp
usr/share/ppp/provider.peer
usr/share/ppp/pap-secrets
usr/share/ppp/chap-secrets
usr/share/ppp/provider.chatscript

