# rootfs/pkgs/pure-ftpd/filelist
etc/pure-ftpd/auth/70pam
etc/pure-ftpd/auth/65unix
etc/pure-ftpd/auth/50pure
etc/pure-ftpd/conf/PureDB
etc/pure-ftpd/conf/NoAnonymous
etc/pure-ftpd/conf/PAMAuthentication
etc/pure-ftpd/conf/MinUID
etc/pure-ftpd/conf/AltLog
etc/pure-ftpd/conf/UnixAuthentication
etc/pure-ftpd/pureftpd.pdb
etc/pure-ftpd/pureftpd.passwd
etc/init.d/pure-ftpd
usr/bin/pure-pw
usr/bin/pure-pwconvert
usr/bin/pure-statsdecode
usr/sbin/pure-ftpd
usr/sbin/pure-ftpd-virtualchroot
usr/sbin/pure-authd
usr/sbin/pure-ftpwho
usr/sbin/pure-quotacheck
