# rootfs/pkgs/cracklibpam/filelist
bin/zcat
var/cache/cracklib/cracklib_dict.hwm
var/cache/cracklib/cracklib_dict.pwd
var/cache/cracklib/cracklib_dict.pwi
lib/security/pam_cracklib.so
etc/logcheck/ignore.d.server/cracklib-runtime
etc/logcheck/ignore.d.paranoid/cracklib-runtime
etc/logcheck/ignore.d.workstation/cracklib-runtime
etc/cracklib/cracklib.conf
etc/cron.daily/cracklib-runtime
usr/bin/uniq
usr/lib/libcrack.so.2
usr/lib/libcrack.so.2.7
usr/sbin/crack_packer
usr/sbin/update-cracklib
usr/sbin/crack_mkdict
usr/sbin/crack_testnum
usr/sbin/crack_teststr
usr/sbin/crack_unpacker
usr/sbin/crack_testlib
